# test-monitor.service
[Unit]
Description=Monitoring process test
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/test_monitor.sh
